# data-clients/mcb-python-client/cron/vgsn/crontab
0,30 * * * * /home/vgsn/Code/MyCityBikes/mycitybikes-server/data-clients/mcb-python-client/cron/vgsn/oslo.sh
